# deploy/solar-dashboard.service
[Unit]
Description=Solar E-Ink Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/solar-eink-dashboard
EnvironmentFile=/home/pi/solar-eink-dashboard/.env.local
ExecStart=/home/pi/solar-eink-dashboard/.venv312/bin/python main.py --production
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
