# deploy/solar-dashboard.service
[Unit]
Description=Solar E-Ink Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%DEPLOY_USER%
WorkingDirectory=%DEPLOY_DIR%
EnvironmentFile=%DEPLOY_DIR%/.env.local
ExecStart=%DEPLOY_DIR%/.venv/bin/python main.py --production
Restart=always
RestartSec=10
WatchdogSec=300
Environment=PYTHONUNBUFFERED=1

# Memory protection: cap the service so OOM kills it (not random system
# processes).  Chromium + Python + numpy typically need ~600-800 MB;
# 1.5 G gives headroom while protecting the remaining ~2.5 G for the OS.
MemoryMax=1536M
MemoryHigh=1280M

[Install]
WantedBy=multi-user.target
