# deploy/solar-dashboard.service
[Unit]
Description=Solar E-Ink Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%DEPLOY_USER%
WorkingDirectory=%DEPLOY_DIR%
EnvironmentFile=%DEPLOY_DIR%/.env.local
ExecStart=%DEPLOY_DIR%/.venv312/bin/python main.py --production
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
